# templates/confluence.service.erb
# Managed by Chef - local modifications will be overwritten
[Unit]
Description=Atlassian Confluence
After=network.target

[Service]
Type=forking
User=<%= @user %>
Group=<%= @group %>
Environment=CONFLUENCE_HOME=<%= @home_path %>
PIDFile=<%= @install_path %>/work/catalina.pid
ExecStart=<%= @install_path %>/bin/start-confluence.sh
ExecStop=<%= @install_path %>/bin/stop-confluence.sh
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
